# quicktill-nfc-bridge@.service
[Unit]
Description=NFC reader to quicktill bridge on %i

[Service]
ExecStart=/usr/sbin/quicktill-nfc-bridge %i
